# fula/fula.service
[Unit]
Description=%i service with docker compose
PartOf=docker.service
After=docker.service

[Service]
Type=simple
WorkingDirectory=/usr/bin/fula/
ExecStart=/usr/bin/fula/fula.sh start
ExecStop=/usr/bin/fula/fula.sh stop
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
